# corescope/rootfs/etc/services.d/corescope-server/finish
#!/usr/bin/with-contenv bashio
exit_code="${1:-unknown}"
term_signal="${2:-none}"
bashio::log.info "CoreScope web server service stopped. exit_code=${exit_code} signal=${term_signal}"
