# corescope/rootfs/etc/services.d/corescope/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

OPTIONS_PATH="/data/options.json"

disable_caddy="$(jq -r '.disable_caddy // true' "${OPTIONS_PATH}")"

export DB_PATH="/app/data/meshcore.db"

if [ "${disable_caddy}" = "true" ]; then
  port="80"
  bashio::log.info "Starting CoreScope web server on port ${port}."
else
  port="3000"
  bashio::log.info "Starting CoreScope web server on port ${port} behind Caddy."
fi

exec /app/corescope-server -config-dir /app -db /app/data/meshcore.db -public /app/public -port "${port}"
