# ingress-lan-proxy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh
set -eu

exec nginx -g "daemon off;"
