# corescope/rootfs/etc/services.d/corescope-caddy/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

OPTIONS_PATH="/data/options.json"

disable_caddy="$(jq -r '.disable_caddy // true' "${OPTIONS_PATH}")"

if [ "${disable_caddy}" = "true" ]; then
  bashio::log.info "Caddy is disabled; CoreScope will serve port 80 directly."
  exec sleep 2147483647
fi

bashio::log.info "Starting Caddy reverse proxy on port 80."
exec /usr/sbin/caddy run --config /config/caddy/Caddyfile
